# packaging/bolt.service
[Unit]
Description=Bolt Download Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=%h/.local/bin/bolt
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/bolt %h/.local/share/bolt %h/Downloads
PrivateTmp=true

[Install]
WantedBy=default.target
